# databricks/apps/databricks-agent-proxy/src/databricks_agent_proxy/templates/systemd.service
[Unit]
Description=Databricks Agent Proxy
After=network.target

[Service]
Type=simple
ExecStart=${uvx_path} databricks-agent-proxy start --gateway-url ${gateway_url} --port ${port}
Restart=always
RestartSec=5
Environment=HOME=${home}

[Install]
WantedBy=default.target
